# config/dbgserver.service.in
[Unit]
Description=
Documentation=http://elfutils.org/
After=network.target

[Service]
EnvironmentFile=-@sysconfdir@/sysconfig/@program_prefix@dbgserver
User=daemon
CacheDirectory=dbgserver
ExecStart=/bin/sh -c '@prefix@/bin/@program_prefix@dbgserver -d $CACHE_DIRECTORY/dbgserver.sqlite -p '"$DBGSERVER_PORT $DBGSERVER_VERBOSE $DBGSERVER_PATHS"
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
